# com.bluecandle.startup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bluecandle.startup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/apple/Documents/gemini/antigravity/scratch/investorbabu/morning_start.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday at 08:30 AM -->
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>30</integer>
            <key>Weekday</key>
            <integer>1</integer>
        </dict>
        <!-- Tuesday at 08:30 AM -->
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>30</integer>
            <key>Weekday</key>
            <integer>2</integer>
        </dict>
        <!-- Wednesday at 08:30 AM -->
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>30</integer>
            <key>Weekday</key>
            <integer>3</integer>
        </dict>
        <!-- Thursday at 08:30 AM -->
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>30</integer>
            <key>Weekday</key>
            <integer>4</integer>
        </dict>
        <!-- Friday at 08:30 AM -->
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>30</integer>
            <key>Weekday</key>
            <integer>5</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/apple/Documents/gemini/antigravity/scratch/investorbabu/morning_start.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/apple/Documents/gemini/antigravity/scratch/investorbabu/morning_start.err</string>
</dict>
</plist>
